# deploy/launchd/com.soa-web.channels.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SoA-Web channel self-healer. Runs scripts/soa-channels every 60s to keep the
  fleet (com.soa-web.server, :7332) reachable through the best working public
  channel (cloudflare → ngrok → lan), and to write a clobber-proof fleet snapshot
  every pass so nothing is ever lost. It NEVER restarts the daemon or touches a
  tab — only the access layer + a save. Install to ~/Library/LaunchAgents/ and
  `launchctl bootstrap gui/$(id -u) <plist>`. Adjust the absolute paths to the
  install user's home.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.soa-web.channels</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/test/.soa-web/scripts/soa-channels</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>AbandonProcessGroup</key>
    <true/>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/test/.soa-web/logs/channels.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/test/.soa-web/logs/channels.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SOA_WEB_PORT</key>
        <string>7332</string>
    </dict>
</dict>
</plist>
